# cron/.crontab
@daily /Users/rboll/bin/upup >> /var/log/upup.log

